# host/dist/moltbot-ipfs.service
[Unit]
Description=Moltbot IPFS Daemon
After=network.target

[Service]
Type=simple
User=moltbot
ExecStart=/usr/local/bin/ipfs daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
