# host/dist/moltbot-ipfs.service
[Unit]
Description=Moltbot IPFS Daemon
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=moltbot
Environment=HOME=/home/moltbot
Environment=IPFS_PATH=/home/moltbot/.ipfs
ExecStartPre=/usr/local/lib/moltbot/moltbot-ipfs-prestart.sh
ExecStart=/usr/local/bin/ipfs daemon --migrate=true
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
